# ansible/files/supervisord.conf
[program:grunt-watch]
command=grunt watch
autostart=true
autorestart=true
directory={{ root_folder }}
stderr_logfile = /var/log/supervisord-assetic-watch-stderr.log
startretries=900
