# ansible/files/supervisord.conf
[program:assetic]
command=php app/console assetic:dump --watch
autostart=true
autorestart=true
directory={{ root_folder }}
stderr_logfile = /var/log/supervisord-assetic-watch-stderr.log
